# tests/fixtures/curated_packages/systemd_auto_enable_install_hook/curated-systemd-auto-enable.install
post_install() {
  systemctl enable curated-systemd-auto-enable.service
}
