# aurascan/assets/aurascan-incident-maintenance.service
[Unit]
Description=AuraScan weekly incident maintenance scan
Documentation=https://github.com/crizzler/AuraScan
ConditionPathExists=/usr/bin/journalctl

[Service]
Type=oneshot
ExecStart=/usr/bin/aurascan incidents --current-boot --capture-maintenance
User=root
Group=root
UMask=0077
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
PrivateTmp=yes
PrivateNetwork=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectClock=yes
ProtectHostname=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
ProcSubset=pid
ReadWritePaths=/var/lib/aurascan/incidents
StateDirectory=aurascan/incidents
StateDirectoryMode=0755
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=~@clock @cpu-emulation @debug @module @mount @obsolete @privileged @raw-io @reboot @swap
SystemCallErrorNumber=EPERM
IPAddressDeny=any
TimeoutStartSec=120
Nice=15
CPUWeight=10
IOWeight=10
IOSchedulingClass=idle
