# aurascan/assets/aurascan-incident-monitor.service
[Unit]
Description=AuraScan previous-boot incident monitor
Documentation=https://github.com/crizzler/AuraScan
After=systemd-journal-flush.service
Wants=systemd-journal-flush.service
ConditionPathExists=/usr/bin/journalctl
OnSuccess=aurascan-incident-safe-autopilot.service

[Service]
Type=oneshot
ExecStart=/usr/bin/aurascan incidents --last-boot --capture-monitor
User=root
Group=root
UMask=0077
UnsetEnvironment=AURASCAN_AI_KEY AURASCAN_OPENAI_API_KEY AURASCAN_ANTHROPIC_API_KEY AURASCAN_DEEPSEEK_API_KEY AURASCAN_GEMINI_API_KEY AURASCAN_OPENROUTER_API_KEY
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
PrivateTmp=yes
PrivateNetwork=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectClock=yes
ProtectHostname=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
ProcSubset=pid
ReadWritePaths=/var/lib/aurascan/incidents
StateDirectory=aurascan/incidents
StateDirectoryMode=0755
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=~@clock @cpu-emulation @debug @module @mount @obsolete @privileged @raw-io @reboot @swap
SystemCallErrorNumber=EPERM
IPAddressDeny=any
TimeoutStartSec=90
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
