# aurascan/assets/aurascan-incident-safe-autopilot.service
[Unit]
Description=AuraScan offline deterministic incident Safe Autopilot
Documentation=https://github.com/crizzler/AuraScan
ConditionPathExists=/usr/bin/aurascan

[Service]
Type=oneshot
ExecCondition=/usr/bin/aurascan incidents --safe-autopilot-enabled
ExecStart=/usr/bin/aurascan incidents --capture-safe-autopilot
User=root
Group=root
UMask=0077
UnsetEnvironment=AURASCAN_AI_KEY AURASCAN_OPENAI_API_KEY AURASCAN_ANTHROPIC_API_KEY AURASCAN_DEEPSEEK_API_KEY AURASCAN_GEMINI_API_KEY AURASCAN_OPENROUTER_API_KEY
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
PrivateTmp=yes
PrivateNetwork=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectClock=yes
ProtectHostname=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=ptraceable
ProcSubset=pid
ReadWritePaths=/var/lib/aurascan/incidents /var/lib/pacman /etc/pacman.d
StateDirectory=aurascan/incidents
StateDirectoryMode=0755
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=~@clock @cpu-emulation @debug @module @mount @obsolete @privileged @raw-io @reboot @swap
SystemCallErrorNumber=EPERM
IPAddressDeny=any
TimeoutStartSec=180
Nice=10
IOSchedulingClass=idle
